# data/systemd/liri-daemon-power.service.in
[Unit]
Description=Power managent handling
PartOf=liri-daemons.target

[Service]
Type=dbus
ExecStart=@ABSOLUTE_LIBEXECDIR@/liri-daemon --module=power
BusName=io.liri.Daemon.Modules.Power
Restart=on-failure
